# snapshot/20251008T103035Z/systemd/mcp-sse-proxy.service
[Unit]\nDescription=mcp-sse-proxy\nAfter=network.target\n\n[Service]\nType=simple\nWorkingDirectory=/opt/mcp\nEnvironmentFile=/opt/mcp/.env\nExecStart=/usr/bin/node /opt/mcp/dist/index.js\nRestart=always\nRestartSec=2\n\n[Install]\nWantedBy=multi-user.target\n