# snapshot/20251008T103035Z/systemd/stas-gw-smoke.service
[Unit]
Description=STAS GW hourly smoke (GET + DRY-RUN DELETE)

[Service]
Type=oneshot
ExecStart=/usr/local/bin/stas_gw_smoke.sh
